# systemd.service
[Unit]
Description=JellyOS Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=jellyos
Group=jellyos
WorkingDirectory=/opt/jellyos
Environment=NODE_ENV=production
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/etc/jellyos/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jellyos

# Security settings
NoNewPrivileges=true
MemoryDenyWriteExecute=true
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Directory creation and permissions
RuntimeDirectory=jellyos
RuntimeDirectoryMode=0755
LogsDirectory=jellyos
LogsDirectoryMode=0755

[Install]
WantedBy=multi-user.target